# thirdparty/PIGPIO/pigpio.service
[Unit]
Description=PIGPIO daemon
After=network.target auditd.service

[Service]
ExecStart=/usr/local/bin/pigpiod -g
#KillMode=process
Restart=on-failure
#RestartPreventExitStatus=255
Type=exec
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=20

[Install]
WantedBy=multi-user.target
Alias=PIGPIO.service
